# templates/aem.service.j2
# {{ ansible_managed }}

[Unit]
Description={{ aem_instance_name }}
After=multi-user.target
Requires=network.target

[Service]
Type=forking
ExecStart={{ aem_home }}/crx-quickstart/bin/start
ExecStop={{ aem_home }}/crx-quickstart/bin/stop
PIDFile={{ aem_home }}/crx-quickstart/conf/cq.pid
User={{ aem_user }}
LimitNOFILE={{ aem_limit_nofile }}

[Install]
WantedBy=multi-user.target